# airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --noreset --noclear --skip-login --login-program /usr/local/bin/baklava-launcher - ${TERM}
